# conf.d/ed-websockify.conf
[program:websockify]
command=/app/conf.d/ed-websockify.sh
autorestart=true
